# roles/wireguard_recovery/templates/wg-recovery-check.service.j2
[Unit]
Description=WireGuard tunnel health check and automatic recovery
After=network.target wg-quick@{{ wg_recovery_interface }}.service

[Service]
Type=oneshot
ExecStart={{ wg_recovery_script_path }}
